# images/services/backups.conf
[program:backups]
command=/scripts/backups_wrapper.sh
priority=10
autostart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
